# playbooks/templates/2ip-speed.service
[Unit]
Description=2ip-speed
After=network.target

[Service]
NonBlocking=true
ExecStart=/usr/bin/speedtest --email=dev@2ip.ru
ExecReload=/bin/kill -HUP $MAINPID
LimitNOFILE=8192
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target