# playbooks/templates/2ip-speed.socket
[Unit]
Description=2ip socket
PartOf=2ip-speed.service

[Socket]
ListenStream=80
NoDelay=true

[Install]
WantedBy=sockets.target