# playbooks/templates/2ip-speed.service
[Unit]
Description=2ip-speed
After=network.target

[Service]
NonBlocking=true
ExecStart=/usr/bin/speedtest --id={{ uniq_id }}
ExecReload=/bin/kill -HUP $MAINPID
LimitNOFILE=8192
Restart=always
RestartSec=5
AmbientCapabilities=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
