# playbooks/templates/2ip-speed.service
[Unit]
Description=2ip speed platform
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/speedtest
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target